# kata-guest-base/extra-nvidia/etc/systemd/system/kata-agent.service.d/30-nvidia-gpu.conf
# GPU guest only: the kata-agent must not accept container workloads until
# the NVIDIA driver, persistence daemon, CDI spec, and CC ready-state are
# all up (the three nvidia-gpu-* units — the systemd NVRC replacement).
# Requires= makes a GPU bring-up failure fail the agent, and with it the
# sandbox: fail-fast at pod creation, not at first CUDA call.
[Unit]
Requires=nvidia-gpu-ready.service
After=nvidia-gpu-ready.service
