# kata-guest-base/extra/etc/systemd/system/scratch-setup.service
[Unit]
# Backs /run/kata-containers/image (the guest-pull image store) with a
# dm-crypt-encrypted scratch disk, so large workload images unpack to disk
# instead of the RAM tmpfs (which caps image size at guest memory). Runs BEFORE
# kata-agent so the mount is in place before any CreateContainer/guest-pull.
# No-op when the VM has no confai-scratch disk. See extra/usr/local/bin/
# scratch-setup.sh and scripts/kata-qemu-scratch-wrapper.sh.
Description=c8s confidential image-store scratch disk (dm-crypt)
Documentation=https://github.com/confidential-dot-ai/c8s/tree/main/kata-guest-base
DefaultDependencies=no
After=local-fs.target
Before=kata-agent.service
[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/lib/c8s/scratch-setup.sh
# never wedge kata-agent if something hangs
TimeoutStartSec=60
[Install]
WantedBy=kata-containers.target
